# stravolt.service
[Unit]
Description=Stravolt Rails App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/ubuntu/stravoltstravolt
Environment=RAILS_ENV=production
EnvironmentFile=/home/ubuntu/stravoltstravolt/.env
ExecStart=/usr/local/bin/bundle exec rails server -p 3000 -b 127.0.0.1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stravolt

[Install]
WantedBy=multi-user.target
